# config/systemd/user/change-wallpaper.service
[Unit]
Description=Sets a new wallpaper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/scripts/wallpaper

[Install]
WantedBy=default.target
