# etc/profile.d/headphones.sh
export HEADPHONES_MAC="70:26:05:CF:A4:33"
