# etc/profile.d/z_user_config.sh
export EDITOR=vim
export TERMINAL=alacritty

export XDG_CONFIG_HOME="$HOME/.config"
export XDG_CACHE_HOME="$HOME/.cache"
export XDG_DATA_HOME="$HOME/.local/share"

export SSH_AUTH_SOCK="$XDG_RUNTIME_DIR/gnupg/S.gpg-agent.ssh"

export PATH="$HOME/.local/bin/scripts:$HOME/.local/bin:$PATH"
