# etc/profile.d/zz_custom.sh
export EDITOR=vim
export TERMINAL=alacritty

export WINIT_HIDPI_FACTOR=1 # Alacritty hack to correctly support multiple monitors

export XDG_CONFIG_HOME="$HOME/.config"
export XDG_CACHE_HOME="$HOME/.cache"

export SSH_AUTH_SOCK="$XDG_RUNTIME_DIR/gnupg/S.gpg-agent.ssh"

export PATH="$HOME/.local/bin/scripts:$HOME/.local/bin:$PATH"
