# etc/udev/rules.d/10-displays.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/local/bin/configure_displays"
